# config/systemd/user/ags.service
[Unit]
Description=AGS (Aylur's Gtk Shell)
After=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStartPre=-/usr/bin/pkill -x ags
ExecStart=/usr/bin/bash -lc 'exec /usr/bin/ags run --gtk 4 --directory "$HOME/.config/ags"'
Restart=on-abnormal
RestartSec=2
TimeoutStopSec=5
KillMode=mixed

[Install]
WantedBy=default.target
